# routerConfig/files/etc/rc.local
#NODO 10.10.5.5
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
python /root/NodoCliente.py > /dev/null
exit 0

